# unit/NetManager.service
[Unit]
Description=Network manager for oakestra

[Service]
Type=simple
ExecStart=/usr/bin/NetManager -p 6000

[Install]
WantedBy=multi-user.target